# local/com.ticketmonitor.check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ticketmonitor.check</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/gyaici/ticket-monitor/local/run-check.sh</string>
    </array>

    <!-- Vérifie toutes les 120 s (2 min) -->
    <key>StartInterval</key>
    <integer>120</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/gyaici/ticket-monitor/local/check.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/gyaici/ticket-monitor/local/check.log</string>
</dict>
</plist>
